# src/queue_server.service
[Service]
ExecStart=/usr/local/mrt2/src/queue_server.py
User=queueer
Group=queueer
